# deploy/audiostash-scand.service
[Unit]
Description=AudioStash scanner daemon

[Service]
WorkingDirectory=/var/www/audiostash/
ExecStart=/var/virtualenvs/audiostash/bin/python -m audiostash.scand.main
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
